# docker/php/supervisor.dev/dev.conf
[program:queue_work]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
command=php artisan queue:work database --tries=3
numprocs=2
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stderr_logfile=/dev/stdout
stdout_logfile=/dev/stdout
priority=300
user=www-data


